# cfai.service
[Unit]
Description=CF_AI Penetration Testing Framework
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=root
WorkingDirectory=/opt/CF_AI
EnvironmentFile=-/opt/CF_AI/.env
ExecStart=/opt/CF_AI/venv/bin/python3 /opt/CF_AI/cfai_server.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cfai
RuntimeMaxSec=86400

[Install]
WantedBy=multi-user.target
